# debian/postinst
#!/bin/sh
# postinst script for lzfs
#
# see: dh_installdeb(1)
/usr/sbin/update-rc.d zfsload defaults
/etc/init.d/zfsload start

