# debian/postrm
#!/bin/sh
# postrm script for lzfs
#
# see: dh_installdeb(1)

/usr/sbin/update-rc.d -f zfsload remove
/bin/rm -f /etc/init.d/zfsload
