# debian/postinst
#!/bin/sh
# postinst script for lzfs
#
# see: dh_installdeb(1)
/sbin/depmod
/usr/sbin/update-rc.d zfsload defaults

